# appliance/root/etc/systemd/system/incus-spawn-vm.service
[Unit]
Description=incus-spawn VM initialization
After=incus.service
Requires=incus.service
After=systemd-networkd.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/incus-spawn-vm-init
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
